# Bash script/memchecker.timer
[Unit]
Description=Timer for memory checker

[Timer]
OnUnitActiveSec=60